# scripts/deploy-app/resources/supervisor.conf
[program:microblog]
command=/bin/bash /home/deploy/microblog/start.sh
directory=/home/deploy/microblog
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log
